# data/bread-snapshot.timer
[Unit]
Description=Run Bread snapshot hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
